# deploy/supervisor-worker.conf
; ============================================================
; Configuración de Supervisor para BecasLaravel - Queue Worker
; ============================================================
;
; INSTRUCCIONES DE DESPLIEGUE EN PRODUCCIÓN:
;
; 1. Instala Supervisor en el servidor (si no lo tienes):
;       sudo apt-get install supervisor
;
; 2. Copia este archivo al servidor:
;       sudo cp supervisor-worker.conf /etc/supervisor/conf.d/becaslaravel-worker.conf
;
; 3. Actualiza 'command' con la ruta real de tu proyecto en el servidor.
;
; 4. Recarga Supervisor:
;       sudo supervisorctl reread
;       sudo supervisorctl update
;       sudo supervisorctl start becaslaravel-worker:*
;
; 5. Asegúrate de que tu .env en el servidor tenga:
;       QUEUE_CONNECTION=database
;
; 6. Corre las migraciones (incluye la tabla 'jobs'):
;       php artisan migrate --force
;
; ============================================================

[program:becaslaravel-worker]
; *** Cambia esta ruta por la ruta real de tu proyecto en el servidor ***
command=php /var/www/BecasLaravel/artisan queue:work database --sleep=3 --tries=3 --max-time=3600

process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true

; Archivo de log del worker (ajusta la ruta si es necesario)
stdout_logfile=/var/www/BecasLaravel/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

; Usuario que corre el proceso (normalmente www-data en Ubuntu/Debian)
user=www-data

; Tiempo de espera al detener el proceso (segundos)
stopwaitsecs=3600
